# roles/driver/templates/dividat-driver.service.j2
[Unit]
Description=Dividat Driver

[Service]
Type=simple
ExecStart=/home/{{ user }}/bin/dividat-driver
User={{ user }}
Group={{ user }}

[Install]
WantedBy=multi-user.target
